# deploy/systemd/megavpn-migrate.service
[Unit]
Description=RTIS MegaVPN Database Migrations
Documentation=https://github.com/rtis-emc2/megavpn
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/megavpn
EnvironmentFile=/etc/megavpn/megavpn.env
Environment=MEGAVPN_MIGRATIONS_DIR=/opt/megavpn/migrations
ExecStart=/opt/megavpn/bin/megavpn-migrate
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=full
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077
